# backup/pbs/dothomelab-appdata-backup.timer
[Unit]
Description=Daily dothomelab appdata backup

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=15m
Unit=dothomelab-appdata-backup.service

[Install]
WantedBy=timers.target
